# crontab
00 15 * * * /usr/local/bin/python /make_orders.py  >> /foodpass.log 2>&1
