# crontab
00 15 * * * /usr/local/bin/python /makeOrders.py  >> /foodpass.log 2>&1
